# tools/deb/postinst
#!/bin/bash
set -e

CONFIG_FILE="/etc/marmot/marmotd.json"
LEGACY_CONFIG_FILE="/etc/marmot.json"
KEY_DIR="/etc/marmot/keys"
PRIVATE_KEY="${KEY_DIR}/private.key"
PUBLIC_KEY="${KEY_DIR}/public.key"
LEGACY_PUBLIC_KEY="${PRIVATE_KEY}.pub"
PLAYBOOK_SRC_DIR="/usr/local/marmot/gateway-playbooks"
PLAYBOOK_DST_DIR="/var/lib/marmot/ansible-playbooks"
PLAYBOOK_TEMPLATE_DIR="${PLAYBOOK_DST_DIR}/templates"

if [ -f "${LEGACY_CONFIG_FILE}" ]; then
    CONFIG_FILE="${LEGACY_CONFIG_FILE}"
fi

if [ -f "${CONFIG_FILE}" ]; then
    NODE_NAME="$(hostname)"
    if [ -n "${NODE_NAME}" ]; then
        sed -i "s/\"node_name\": \"[^\"]*\"/\"node_name\": \"${NODE_NAME}\"/" "${CONFIG_FILE}"
    fi
fi

mkdir -p "${KEY_DIR}"
chmod 700 "${KEY_DIR}"

KEY_INIT_STATUS="generated"

# Backward compatibility: previous installer created private.key.pub.
if [ -f "${LEGACY_PUBLIC_KEY}" ] && [ ! -f "${PUBLIC_KEY}" ]; then
    mv "${LEGACY_PUBLIC_KEY}" "${PUBLIC_KEY}"
fi

# If any key file already exists under /etc/marmot/keys, keep it as-is.
if find "${KEY_DIR}" -mindepth 1 -maxdepth 1 -type f \( -name '*.key' -o -name '*.pub' \) | grep -q .; then
    KEY_INIT_STATUS="skipped"
else
    if [ ! -f "${PRIVATE_KEY}" ] || [ ! -f "${PUBLIC_KEY}" ]; then
        rm -f "${PRIVATE_KEY}" "${PUBLIC_KEY}" "${LEGACY_PUBLIC_KEY}"
        ssh-keygen -t rsa -b 4096 -f "${PRIVATE_KEY}" -N "" -C "marmot-gateway" >/dev/null

        # ssh-keygen always writes public key as ${PRIVATE_KEY}.pub.
        if [ -f "${LEGACY_PUBLIC_KEY}" ] && [ ! -f "${PUBLIC_KEY}" ]; then
            mv "${LEGACY_PUBLIC_KEY}" "${PUBLIC_KEY}"
        fi
    fi
fi

if [ -f "${PRIVATE_KEY}" ]; then
    chmod 600 "${PRIVATE_KEY}"
fi
if [ -f "${PUBLIC_KEY}" ]; then
    chmod 644 "${PUBLIC_KEY}"
fi

rm -rf "${PLAYBOOK_DST_DIR}"
mkdir -p "${PLAYBOOK_TEMPLATE_DIR}"
if [ -d "${PLAYBOOK_SRC_DIR}" ]; then
    cp -a "${PLAYBOOK_SRC_DIR}/." "${PLAYBOOK_TEMPLATE_DIR}/"
fi

systemctl enable libvirtd.service
systemctl start libvirtd.service

systemctl enable iscsid.service
systemctl start iscsid.service

enable_and_start_if_exists() {
    unit_name="$1"
    if systemctl list-unit-files | grep -q "^${unit_name}\\.service"; then
        systemctl enable "${unit_name}.service" || true
        systemctl start "${unit_name}.service" || true
    fi
}

# OVS runtime
enable_and_start_if_exists openvswitch-switch
enable_and_start_if_exists ovsdb-server
enable_and_start_if_exists ovs-vswitchd

# OVN runtime
enable_and_start_if_exists ovn-central
enable_and_start_if_exists ovn-northd
enable_and_start_if_exists ovn-controller
enable_and_start_if_exists ovn-host

systemctl daemon-reload
systemctl enable marmot.service
systemctl start marmot.service

echo ""
echo "============================================================"
echo " marmot のインストールが完了しました"
echo "============================================================"
if [ "${KEY_INIT_STATUS}" = "skipped" ]; then
    echo " Gateway playbook を初期化しました (既存の SSH 鍵を検出したため鍵初期化はスキップ)"
else
    echo " Gateway playbook と SSH 鍵を初期化しました"
fi
echo " 自動的に VM起動イメージをダウンロードしています..."
echo " この処理は環境によって数分かかることがあります。"
echo " mactl get images でダウンロード状況を確認できます。"
echo "============================================================"

exit 0
